# deploy/linux/hechao-launcher-api.service
[Unit]
Description=Hechao Launcher API
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=simple
User=hechao-api
Group=hechao-api
WorkingDirectory=/opt/hechao-launcher-api/current
ExecStart=/opt/hechao-launcher-api/current/Hechao.Api
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://127.0.0.1:8090
EnvironmentFile=/etc/hechao-launcher-api/environment
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
