# deploy/linux/database/hechao-launcher-db-backup.service
[Unit]
Description=Back up the Hechao launcher PostgreSQL database
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hechao-launcher-db-backup
User=root
Group=root
UMask=0077
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
